# files/camera-stream.socket
[Unit]
Description=Camera Stream Socket

[Socket]
ListenStream=5000
Accept=yes

[Install]
WantedBy=sockets.target
